# rootfs/etc/crontabs/update_tracker
SHELL=/bin/bash
0	7	*	*	*	sleep $((RANDOM % 1800)); bash /config/aria2/script/tracker.sh /config/aria2/aria2.conf RPC 2>&1 | tee /config/aria2/tracker.log
