# rootfs/etc/services.d/nmbd/finish
#!/usr/bin/with-contenv bash

echo "SERVICE ENDED: nmbd-service"
#s6-svscanctl -t /var/run/s6/services