# rootfs/etc/services.d/nginx/run
#!/usr/bin/with-contenv bash

echo "STARTED SERVICE: nginx-service"

exec nginx -c /config/nginx/nginx.conf -g 'daemon off;'