# rootfs/etc/services.d/aria2/finish
#!/usr/bin/with-contenv bash

echo "SERVICE ENDED: aria2-service"
#s6-svscanctl -t /var/run/s6/services
