# rootfs/etc/services.d/aria2/run
#!/usr/bin/with-contenv bash

echo "STARTED SERVICE: aria2-service"
umask ${UMASK_SET:-022}
exec s6-setuidgid p3terx aria2c \
    --conf-path=/config/aria2/aria2.conf
