# rootfs/etc/services.d/alist/finish
#!/usr/bin/with-contenv bash

echo "SERVICE ENDED: alist-service"
#s6-svscanctl -t /var/run/s6/services
