# rootfs/etc/services.d/alist/run
#!/usr/bin/with-contenv bash

echo "SERVICE STARTED: alist-service"
umask ${UMASK_SET:-022}
exec s6-setuidgid p3terx alist server --data /config/alist --no-prefix
